# deploy/services/mower-ros2-bridge.service
[Unit]
Description=Mower-Roboto ROS 2 Serial Bridge
After=network.target mower-backend.service
Requires=network.target

[Service]
Type=simple
User=pi
Environment="ROS_DOMAIN_ID=0"
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/pi/mower-roboto/ros2/install/setup.bash && ros2 launch mower_bringup serial_bridge.launch.py'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
